# bitrix/crontab/crontab.cfg
0 */24 * * * /usr/local/php/bin/php -f /var/www/html/test.maxclean.help/bitrix/php_interface/include/catalog_export/cron_frame.php 2 >/var/www/html/test.maxclean.help/bitrix/php_interface/include/catalog_export/logs/2.txt
